# creation/templates/gwms-renew-proxies.cron
23 * * * * [ -f /var/lock/subsys/gwms-renew-proxies ] && /sbin/renew_proxies