# creation/templates/gwms-frontend.service
[Unit]
Description=GWMS Frontend Service
Requires=condor.service
After=condor.service
Documentation="http://glideinwms.fnal.gov/doc.prd/frontend/index.html"

[Service]
Type=forking
ExecStart=/usr/sbin/gwms-frontend start
ExecStop=/usr/sbin/gwms-frontend stop
ExecReload=/usr/sbin/gwms-frontend reload  $MAINPID

[Install]
WantedBy=multi-user.target
