# creation/templates/gwms-renew-proxies.cron
23 * * * * root [ -f /var/lock/subsys/gwms-renew-proxies ] && sleep $((RANDOM % 600)) && /usr/libexec/gwms_renew_proxies
